# server/supervisord/daphne.ini
[program:daphne]
command=/usr/local/python3/bin/daphne -b 0.0.0.0 -p 9001 -u /tmp/daphne.sock django_proj.asgi:application
directory=/www/wwwroot/api.freaks.group/django_proj/
stdout_logfile=/www/wwwroot/api.freaks.group/server/supervisord/daphne_out.log
stderr_logfile=/www/wwwroot/api.freaks.group/server/supervisord/daphne_err.log
autostart=true
autorestart=true
user=root
startsecs=10